# services/transfer_switch/run
#!/bin/sh
echo "*** starting Transfer Switch Monitor ***"
exec 2>&1
exec /data/TransferSwitch/transfer_switch.py
